# ops/systemd/shengfukung-wenfu-puma.service
# ops/systemd/shengfukung-wenfu.service
#
# Systemd unit template for the Puma web server.
#
# NOTE:
# - Replace the placeholder slug (遶ｹ蜊鈴式閨也ｦ丞ｮｮ) with `AppConstants::Project::SLUG` if you rename the stack.
# - Keep `EnvironmentFile=/etc/default/<slug>-env` in sync with the supplied `.env` helpers.

[Unit]
Description=遶ｹ蜊鈴式閨也ｦ丞ｮｮ Puma HTTP Server
After=network.target

[Service]
Type=simple
User=jimmy1768_user
WorkingDirectory=/home/jimmy1768_user/Projects/shengfukung-wenfu/rails
EnvironmentFile=/etc/default/shengfukung-wenfu-env
ExecStart=/usr/bin/bash -lc "/home/jimmy1768_user/.rbenv/bin/rbenv exec bundle exec puma -C config/puma.rb"
Restart=always

[Install]
WantedBy=multi-user.target
